# docker/supervisor/crontab
*/30 * * * * /usr/local/bin/php /var/www/form/artisan app:tasks:refresh
